# supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=/usr/bin/python /app/app.py
